# com.southmoon.cameras.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.southmoon.cameras</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/bin/zsh</string>
        <string>/Users/pennrobotics/Documents/GitHub/Southmoon/start_cameras.sh</string>
    </array>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>
    
    <key>StandardOutPath</key>
    <string>/Users/pennrobotics/Library/Logs/Southmoon/camera_stdout.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/pennrobotics/Library/Logs/Southmoon/camera_stderr.log</string>
    
    <key>WorkingDirectory</key>
    <string>/Users/pennrobotics/Documents/GitHub/Southmoon</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
</dict>
</plist>
